# deploy/systemd/cook-system-worker.service
[Unit]
Description=Cook System TaskIQ Worker
After=network.target redis-server.service cook-system-api.service
Wants=redis-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/user/cook-system
EnvironmentFile=/home/user/cook-system/.env
ExecStart=/usr/local/bin/uv run taskiq worker backend.broker:broker backend.tasks
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
